# src/supervisord.conf
[supervisord]
nodaemon=true

[program:sendmail]
command=/bin/bash -c "exec /etc/init.d/sendmail start"

[program:php-fpm]
command=/bin/bash -c "exec /usr/local/sbin/php-fpm"
